# deploy/supervisor/nozu-workers.conf
[program:nozu-import]
process_name=%(program_name)s
command=/bin/bash -lc "cd /var/www/nozu.me && php artisan nozume:import-queue && exec php artisan queue:work database --queue=imports --sleep=2 --timeout=300 --max-jobs=500"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/nozu.me/storage/logs/import-worker.log
stopwaitsecs=360

[program:nozu-scanner]
process_name=%(program_name)s
command=/bin/bash -lc "cd /var/www/nozu.me && php artisan nozume:sync-resume && exec php artisan queue:work database --queue=scanner --sleep=2 --timeout=180 --max-jobs=300"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/nozu.me/storage/logs/scanner-worker.log
stopwaitsecs=240
